# SOURCES/s3_auction_upload.timer
[Unit]
Description=

[Timer]
OnCalendar=hourly
Persistent=false
Unit=s3_auction_upload.service

[Install]
WantedBy=basic.target
